# src/gnome-remote-desktop.service.in
[Unit]
Description=GNOME Remote Desktop

[Service]
Type=dbus
BusName=org.gnome.RemoteDesktop
ExecStart=@libexecdir@/gnome-remote-desktop-daemon
Restart=on-failure

[Install]
WantedBy=gnome-session.target
